# deploy/checkin-bot.service
[Unit]
Description=Voice health check-in Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/opt/checkin
EnvironmentFile=/opt/checkin/.env
ExecStart=/opt/checkin/.venv/bin/python -m src.main
Restart=on-failure
RestartSec=10s
StandardOutput=append:/var/log/checkin/stdout.log
StandardError=append:/var/log/checkin/stderr.log

[Install]
WantedBy=multi-user.target
